# etc/crontab
# min	hour	mday	month	wday
*/6	*	*	*	*	$HOME/lib/exec/update-dns
*/30	*	*	*	*	pgrep -x -U "$(id -u)" Xorg > /dev/null 2>&1 && env "PATH=$HOME/bin:$PATH" remote-display randomise-backgrounds -d
