# etc/crontab
*/6	*	*	*	*	sh $HOME/lib/update-dns.sh
*/30	*	*	*	*	pgrep -x Xorg >/dev/null 2>&1 && env PATH=$HOME/bin:$PATH remote-display randomise-wallpaper -s
