# etc/sv/scron/run
#!/bin/sh
exec 2>&1
exec scrond -n -f "$HOME/etc/crontab"
