# etc/sv/keybase/run
#!/bin/sh
exec 2>&1
export XDG_RUNTIME_DIR="${XDG_RUNTIME_DIR:-"/run/user/$( id -u )"}"
keybase ctl init
keybase ctl start
exec chpst -b keybase pause
